# images/admin-php/files/etc/cron.d/os2display-cron
# every minute of every day for the rest of your life, run this script as the
# `app` user.
* * * * * root /usr/local/sbin/run_cron.sh $ENV
